# deploy/systemd/e87canbus.service
[Unit]
Description=E87 CAN Bus Bench App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
PermissionsStartOnly=true
WorkingDirectory=/opt/e87canbus
ExecStartPre=/usr/local/bin/e87canbus-can-up
ExecStart=/home/pi/.local/bin/uv run e87canbus-bench-pingpong --interface can0
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
