# deploy/systemd/e87canbus-kiosk.service
[Unit]
Description=E87 kiosk display (cage)
After=e87canbus-controller.service systemd-logind.service

[Service]
Type=simple
User=KIOSK_USER_PLACEHOLDER
PAMName=login
Environment=XDG_SEAT=seat0
Environment=XDG_VTNR=7
TTYPath=/dev/tty7
StandardInput=tty
ExecStartPre=-/usr/bin/plymouth quit --retain-splash
ExecStart=/usr/bin/cage -s -- /usr/local/bin/e87canbus-kiosk
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
